# systemd/solarcontrol-haslach.service
[Unit]
Description=SolarControl Controller (Haslach)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__SERVICE_USER__
WorkingDirectory=__WORKING_DIR__
ExecStart=__EXEC_START__

Restart=on-failure
RestartSec=5

Environment=PYTHONUNBUFFERED=1

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
